# scripts/deployment/supervisor-boatrace-standard-evaluation.ini
[program:boatrace-standard-evaluation]
command=/workspace/boat/scripts/run_standardized_evaluation_cycle.sh
directory=/workspace/boat
priority=60
autostart=true
autorestart=unexpected
startsecs=5
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/workspace/boat/logs/standardized_365d_v2/supervisor.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/workspace/boat/logs/standardized_365d_v2/supervisor-error.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/workspace/boat/src",PGPASSFILE="/workspace/postgresql/conf/databases/boatrace.pgpass",BOATRACE_DB="postgresql://boatrace_app@127.0.0.1:5432/boatrace",BOATRACE_EVAL_CYCLE_SECONDS="21600",BOATRACE_EVAL_RETRY_SECONDS="900",BOATRACE_MODEL_DIR="/workspace/boat/data/models",BOATRACE_LOG_DIR="/workspace/boat/logs"
