# scripts/deployment/supervisor-container-base.ini
[program:container-base]
command=/start.sh
directory=/
priority=10
autostart=true
autorestart=false
startsecs=3
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/workspace/service-manager/log/container-base.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
redirect_stderr=true
