# scripts/deployment/supervisor-boatrace-shadow-provisional.ini
[program:boatrace-odds-shadow-provisional]
command=/workspace/boat/scripts/deployment/run-boatrace-shadow-cycle.sh
directory=/workspace/boat
priority=45
autostart=true
autorestart=true
startsecs=10
startretries=100
stopwaitsecs=60
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/workspace/boat/logs/runtime/realtime-odds-shadow-provisional.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
redirect_stderr=true
environment=BOATRACE_SHADOW_LOCK="/workspace/boat/run/realtime-odds-shadow-provisional.lock",BOATRACE_SHADOW_MODEL="/workspace/boat/data/models/realtime_odds_shadow_t5_safe_candidate.joblib",BOATRACE_SHADOW_BACKTEST="/workspace/boat/data/models/realtime_odds_shadow_t5_safe_candidate_backtest.json",BOATRACE_SHADOW_STATE="/workspace/boat/data/models/realtime_odds_shadow_t5_safe_candidate_state.json",BOATRACE_SHADOW_TARGET_RACES="450",BOATRACE_SHADOW_MIN_TRAIN_RACES="400",BOATRACE_SHADOW_MIN_NEW_RACES="50",BOATRACE_SHADOW_FOLDS="3",BOATRACE_SHADOW_INTERVAL="1800"
