# scripts/deployment/supervisor-boatrace-work-tracking-sync.ini
[program:boatrace-work-tracking-sync]
command=/workspace/boat/.venv/bin/python -m boatrace_ai.work_tracking_sync --db "host=127.0.0.1 port=5432 dbname=boatrace user=boatrace_app" --repo ryo100794/boat --direction both --apply --interval-seconds 900
directory=/workspace/boat
priority=65
autostart=true
autorestart=true
startsecs=5
startretries=20
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/workspace/boat/logs/work-tracking-sync.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
redirect_stderr=true
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/workspace/boat/src",PGPASSFILE="/workspace/postgresql/conf/databases/boatrace.pgpass"
