# scripts/deployment/supervisor-postgresql.ini
[program:postgresql]
command=/workspace/shared-services/bin/run-postgresql.sh
directory=/workspace/postgresql
priority=20
autostart=true
autorestart=unexpected
exitcodes=0
startsecs=5
startretries=10
stopwaitsecs=120
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/workspace/postgresql/log/supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
redirect_stderr=true
