# scripts/deployment/supervisor-boatrace-dashboard.ini
[program:boatrace-dashboard]
command=/workspace/shared-services/bin/run-boatrace-dashboard.sh
directory=/workspace/boat
priority=40
autostart=true
autorestart=true
startsecs=5
startretries=20
stopasgroup=true
killasgroup=true
stdout_logfile=/workspace/boat/data/logs/web-dashboard.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/workspace/boat/data/logs/web-dashboard-error.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3
environment=PYTHONUNBUFFERED="1"
