# scripts/deployment/supervisor-boatrace-conditional-market-shadow.ini
[program:boatrace-conditional-market-shadow]
command=/workspace/boat/scripts/deployment/run-boatrace-conditional-market-shadow.sh
directory=/workspace/boat
priority=49
autostart=true
autorestart=true
startsecs=10
startretries=100
stopwaitsecs=60
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/workspace/boat/logs/runtime/conditional-market-shadow.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
redirect_stderr=true
