# scripts/deployment/supervisor-boatrace-market-blend-shadow.ini
[program:boatrace-market-blend-shadow]
command=/workspace/boat/scripts/deployment/run-boatrace-market-blend-shadow.sh
directory=/workspace/boat
priority=48
autostart=true
autorestart=true
startsecs=10
startretries=100
stopwaitsecs=60
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/workspace/boat/logs/runtime/market-calibrated-blend-shadow.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
redirect_stderr=true
