# scripts/deployment/supervisor-boatrace-collector.ini
[program:boatrace-collector]
command=/workspace/shared-services/bin/run-boatrace-collector.sh
directory=/workspace/boat
priority=30
autostart=true
autorestart=true
startsecs=10
startretries=100
stopwaitsecs=60
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/workspace/boat/logs/runtime/postgresql-collector.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
redirect_stderr=true
