# scripts/deployment/supervisor-boatrace-raw-archive.ini
[program:boatrace-raw-archive]
command=/workspace/shared-services/bin/run-boatrace-raw-archive.sh
directory=/workspace/boat
priority=50
autostart=true
autorestart=true
startsecs=5
startretries=20
stopasgroup=true
killasgroup=true
stdout_logfile=/workspace/boat/data/logs/raw-archive.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/workspace/boat/data/logs/raw-archive-error.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3
environment=PYTHONUNBUFFERED="1",RCLONE_CONFIG="/workspace/google-drive/rclone.conf",BOATRACE_RAW_ARCHIVE_INTERVAL="600",BOATRACE_RAW_ARCHIVE_MIN_AGE_MINUTES="30",BOATRACE_RAW_ARCHIVE_BATCH_FILES="5000"
