# scripts/deployment/supervisor-boatrace-evaluation-runner.ini
[program:boatrace-evaluation-runner]
command=/workspace/boat/.venv/bin/python -m boatrace_ai.evaluation_queue run --db "host=127.0.0.1 port=5432 dbname=boatrace user=boatrace_app" --app-root /workspace/boat --python /workspace/boat/.venv/bin/python --seed-defaults --poll-seconds 5 --seed-interval 3600 --stale-minutes 180 --vm-limit-gib 0 --nice 10
directory=/workspace/boat
process_name=%(program_name)s-%(process_num)02d
numprocs=4
priority=61
autostart=true
autorestart=true
startsecs=5
startretries=20
stopwaitsecs=60
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/workspace/boat/logs/evaluation_queue/runner-%(process_num)02d.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
redirect_stderr=true
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/workspace/boat/src",PGPASSFILE="/workspace/postgresql/conf/databases/boatrace.pgpass",OMP_NUM_THREADS="4",OPENBLAS_NUM_THREADS="4",MKL_NUM_THREADS="4"
