# scripts/deployment/supervisor-boatrace-market-predictor.ini
[program:boatrace-market-predictor]
command=/workspace/boat/scripts/deployment/run-boatrace-market-predictor.sh
directory=/workspace/boat
autostart=true
autorestart=true
startsecs=5
startretries=20
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/workspace/boat/logs/runtime/market-predictor.log
stderr_logfile=/workspace/boat/logs/runtime/market-predictor-error.log
stdout_logfile_maxbytes=20MB
stderr_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile_backups=3
environment=PYTHONUNBUFFERED="1"
